# deploy/social-drip/xyz.sageart.social-drip.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  SAGE Social activity drip — keeps testnet.sageart.xyz alive until launch.
  Install:
    cp deploy/social-drip/xyz.sageart.social-drip.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/xyz.sageart.social-drip.plist
  Stop (at launch, or any time):
    launchctl unload ~/Library/LaunchAgents/xyz.sageart.social-drip.plist
  Runs every 240s (4 min); the script self-jitters up to ~30s so ticks
  don't land on a fixed clock edge.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>xyz.sageart.social-drip</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/dannyelcristyenciobanica/Downloads/SAGE/deploy/social-drip/run-drip.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>240</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/dannyelcristyenciobanica/Downloads/SAGE/deploy/social-drip/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/dannyelcristyenciobanica/Downloads/SAGE/deploy/social-drip/launchd.err</string>
</dict>
</plist>
